# debian/opi-access.service
[Unit]
Description=OPI access renews signed certificates
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/share/opi-access/opi-access-boot.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
